# sock5.sh
#!/bin/sh

# é»˜è®¤é…ç½®ï¼ˆå¦‚æœªä¼ å…¥çŽ¯å¢ƒå˜é‡åˆ™ä½¿ç”¨è¿™äº›ï¼‰
PORT=${PORT:-10808}
USERNAME=${USERNAME:-"user"}
PASSWORD=${PASSWORD:-"pass"}

# æ£€æŸ¥ç³»ç»Ÿç±»åž‹
if [ -f /etc/alpine-release ]; then
  echo "[INFO] æ£€æµ‹åˆ° Alpine ç³»ç»Ÿï¼Œå®‰è£…ä¾èµ–ä¸­..."
  apk update && apk add curl tar
else
  echo "[INFO] æ£€æµ‹åˆ°éž Alpine ç³»ç»Ÿï¼Œå°è¯•å®‰è£… curl å’Œ tar..."
  command -v apt && apt update && apt install -y curl tar || true
  command -v yum && yum install -y curl tar || true
fi

# å®‰è£… sing-box
mkdir -p /etc/sing-box
cd /etc/sing-box || exit 1
echo "[INFO] ä¸‹è½½ sing-box ä¸­..."
curl -L -o sing-box.tar.gz https://github.com/SagerNet/sing-box/releases/latest/download/sing-box-linux-amd64.tar.gz
tar -xzf sing-box.tar.gz
cp sing-box /usr/local/bin/
chmod +x /usr/local/bin/sing-box

# å†™å…¥é…ç½®æ–‡ä»¶
cat > /etc/sing-box/config.json <<EOF
{
  "log": {
    "level": "info"
  },
  "inbounds": [
    {
      "type": "socks",
      "listen": "0.0.0.0",
      "listen_port": $PORT,
      "users": [
        {
          "username": "$USERNAME",
          "password": "$PASSWORD"
        }
      ]
    }
  ],
  "outbounds": [
    {
      "type": "direct"
    }
  ]
}
EOF

# å¯åŠ¨æœåŠ¡
if [ -f /etc/alpine-release ]; then
  echo "[INFO] ä½¿ç”¨ OpenRC å¯åŠ¨ sing-box..."
  cat > /etc/init.d/sing-box <<'RC'
#!/sbin/openrc-run
description="sing-box socks5 service"
command=/usr/local/bin/sing-box
command_args="run -c /etc/sing-box/config.json"
RC
  chmod +x /etc/init.d/sing-box
  rc-update add sing-box
  rc-service sing-box restart
else
  echo "[INFO] ä½¿ç”¨ systemd å¯åŠ¨ sing-box..."
  cat > /etc/systemd/system/sing-box.service <<EOF
[Unit]
Description=Sing-box Socks5 Service
After=network.target

[Service]
ExecStart=/usr/local/bin/sing-box run -c /etc/sing-box/config.json
Restart=always

[Install]
WantedBy=multi-user.target
EOF

  systemctl daemon-reexec
  systemctl daemon-reload
  systemctl enable sing-box
  systemctl restart sing-box
fi

# è¾“å‡ºæç¤º
echo "âœ… Socks5 å¯åŠ¨æˆåŠŸ"
echo "ðŸ“Œ IPåœ°å€ï¼šä½ çš„æœåŠ¡å™¨å…¬ç½‘ IP"
echo "ðŸ“Œ ç«¯å£ï¼š$PORT"
echo "ðŸ“Œ ç”¨æˆ·åï¼š$USERNAME"
echo "ðŸ“Œ å¯†ç ï¼š$PASSWORD"